# deploy/signglove.service
[Unit]
Description=SignGlove Data Collection Server
Documentation=https://github.com/SignGlove/server
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=forking
User=signglove
Group=signglove
WorkingDirectory=/opt/signglove

# 환경변수 파일 로드
EnvironmentFile=/opt/signglove/config/.env

# Docker Compose로 서버 시작
ExecStartPre=/usr/bin/docker-compose -f /opt/signglove/docker-compose.yml down
ExecStart=/usr/bin/docker-compose -f /opt/signglove/docker-compose.yml up -d
ExecStop=/usr/bin/docker-compose -f /opt/signglove/docker-compose.yml down
ExecReload=/usr/bin/docker-compose -f /opt/signglove/docker-compose.yml restart

# 서비스 재시작 정책
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=120

# 리소스 제한
LimitNOFILE=65536
LimitNPROC=4096

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/signglove/data /opt/signglove/backup /var/log/signglove

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=signglove

[Install]
WantedBy=multi-user.target 